# conf/funkwhale.target
[Unit]
Description=Funkwhale (__APP__)
Wants=__APP__-server.service __APP__-worker.service __APP__-beat.service
